# ansible/roles/saleor/templates/uwsgi_supervisor.conf
[program:uwsgi]
command=uwsgi --emperor /etc/uwsgi/vassals --uid www-data --gid www-data
user=www-data
priority=400
autostart=true
autorestart=true

